# root/etc/s6-overlay/s6-rc.d/svc-act_runner/run
#!/command/with-contenv bash
# shellcheck shell=bash

export CONFIG_FILE="/data/config.yaml"

exec \
  /usr/local/bin/run.sh